# usecases/moviering/moviering.service
[Unit]
Description=Movie Ring — Memoreei movie mentions web app
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/fi/.openclaw/workspace-elliot/projects/memoreei/usecases/moviering
ExecStart=/home/fi/.openclaw/workspace-elliot/projects/memoreei/.venv/bin/python app.py
Restart=on-failure
RestartSec=5
Environment=MOVIERING_PORT=5050
EnvironmentFile=/home/fi/.openclaw/workspace-elliot/projects/memoreei/.env

[Install]
WantedBy=default.target
